# board/excito/common/rootfs-step2-overlay/etc/network/if-up.d/ntpdate
#!/bin/sh

if [ "$IFACE" != "lo" ] ; then
	/usr/bin/ntpdate 0.pool.ntp.org
fi
