# board/excito/common/rootfs-step2-overlay/etc/network/if-up.d/rdate.sh
#!/bin/sh

if [ "$IFACE" != "lo" ] ; then
	/usr/sbin/rdate time.nist.gov
fi
