# board/excito/common/rootfs-step2-overlay/etc/profile.d/bash-prompt.sh
if [ "$BASH" ]; then
    export PS1="[\u@\h \W]\\$ "
fi
